# systemd/overrides/officium-api.service.d/ordering.conf
[Unit]
After=network.target llm-stack.service
Wants=llm-stack.service
